# cron.target
[Unit]
Description=SSSD/cron barrier
After=sssd.service network-online.target
Before=crond.service

[Install]
WantedBy=multi-user.target
