# httpd.target
[Unit]
Description=DNS/HTTPD barrier
After=named.service dnsmasq.service
Before=httpd.service

[Install]
WantedBy=multi-user.target
