# keepalived.target
[Unit]
Description=OpenVPN/keepalived barrier
After=openvpn.service
Before=keepalived.service

[Install]
WantedBy=multi-user.target
