# named.target
[Unit]
Description=Network-online / named daemon barrier
After=network-online.target
Before=named.service

[Install]
WantedBy=multi-user.target
