# sssd.target
[Unit]
Description=DNS/SSSD barrier
After=named.service dnsmasq.service
Before=sssd.service

[Install]
WantedBy=multi-user.target
