# zabbix-agent.target
[Unit]
Description=SSSD/Zabbix barrier
After=sssd.service
Before=zabbix-agent.service

[Install]
WantedBy=multi-user.target
